# backend/ops/systemd/lpc-incident-runner.service.example
[Unit]
Description=LPC incident runner
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/opt/lets-paraconnect/backend
Environment=NODE_ENV=production
EnvironmentFile=/etc/letsparaconnect/backend.env
ExecStart=/usr/bin/npm run incident:runner
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=960
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
